# rc.local
#!/bin/sh -e
/opt/securitycam.py
exit 0
